# itamae/cookbooks/letsencrypt_renew/templates/etc/systemd/system/certbot.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=<%= node[:letsencrypt_renew][:certbot_auto_path] %> renew <%= @options %>
